# systemd/user/songbook-backup.timer
[Unit]
Description=Run Songbook-OC database backup daily

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
RandomizedDelaySec=30m
Unit=songbook-backup.service

[Install]
WantedBy=timers.target
