# systemd/user/songbook-backup.service
[Unit]
Description=Songbook-OC SQLite backup

[Service]
Type=oneshot
WorkingDirectory=%h/songbook-oc
ExecStart=%h/songbook-oc/bin/songbook-backup
